# tgbot_template.service
[Unit]
Description=CSUST TG Bot

[Service]
Type=forking
ExecStart=/usr/bin/python3 {}

[Install]
WantedBy=multi-user.target
